# .config/systemd/user/nwm.service
[Unit]
Description=nwm

[Service]
ExecStart=/usr/bin/nwm
Environment=DISPLAY=:0
Restart=on-failure

[Install]
WantedBy=wm.target
